# profile.d/apps
#!/bin/bash

# Defined prefered applications
export TERMINAL=urxvtc
export PAGER=less
export EDITOR=vim
export BROWSER=surf

