# profile.d/keychain
#!/bin/bash

if (command -v keychain &>/dev/null); then
  eval "$(keychain --eval --noask ~/.ssh/{id_ed25519,id_rsa_native})"
fi
