# profile.d/keychain
#!/bin/bash

if $(which keychain &> /dev/null)
then
  eval $(keychain --agents ssh,gpg --eval)
fi
