# profile.d/lang
#!/bin/bash

# Define prefered LANG
export LANG=fr_FR.utf8
