# profile.d/mail
#!/bin/bash

# Mail
export MAILPATH=~/Mails/
