# profile.d/tsocks
#!/bin/bash

# Proxy
if [ -x /usr/bin/tsocks ] then
  export TSOCKS_CONF_FILE=~/.config/tsocks/tsocks.conf
fi
