# profile.d/tsocks
#!/bin/bash

# Proxy
export TSOCKS_CONF_FILE=~/.config/tsocks/tsocks.conf
