# profile.d/xdg
#!/bin/bash

# Defined prefered applications
export XDG_RUNTIME_DIR=$HOME/.cache
